# cron/root.cron
#Root cron 
#restart n8n and update
5 0 * * * /root/bin/n8n_server/bin/start_n8n.sh
#
# OS update
15 0 * * * /root/bin/n8n_server/bin/update.sh
#
# end